# images/bootstrap/overlay/etc/systemd/system/load-images.service
[Unit]
Description=Load pre-cached images
ConditionPathExistsGlob=/usr/lib/container-images/*.tar
Before=microshift.service
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/load-images.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target